# dist/prise.service
[Unit]
Description=Prise terminal multiplexer server
Documentation=https://prise.sh

[Service]
Type=simple
+ExecStart=%h/.local/bin/prise serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
